# buildroot-external/generic-rootfs/board/generic-rootfs/overlays/debug/etc/profile.d/debug.sh
# Debug profile (sourced from /etc/profile.d when debug overlay is used)
export PS1='\u@\h:\w\$ '
alias ll='ls -la'
